# private_Library/private_LaunchAgents/com.filipeestacio.secondbrain-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.filipeestacio.secondbrain-sync</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/filipe/.local/bin/secondbrain-sync</string>
    </array>
    <key>StartInterval</key>
    <integer>300</integer>
    <key>StandardOutPath</key>
    <string>/Users/filipe/Library/Logs/secondbrain-sync.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/filipe/Library/Logs/secondbrain-sync.log</string>
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
